# systemd-user/openclaw-gateway.service.d/auto-fix.conf
[Unit]
# If the gateway enters failed state (after hitting start limits), trigger the fixer.
OnFailure=openclaw-fix.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Restart=always
